# systemd/milamont-observatory.timer
[Unit]
Description=Run the Milamont observatory build once a day

[Timer]
OnCalendar=*-*-* 06:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
